# supervisor/cron-supervisor.conf
[program:cron]
command=service cron start && cron -f
autostart=true
autorestart=false